# packages/lbl-qubic/lbl_qubic-25.8.0.tar.gz/lbl_qubic-25.8.0/scripts/qubic_job_server.service
[Unit]
Description=Pinto Job Server

[Service]
ExecStart=/home/qubic/anaconda3/envs/qubic/bin/python /home/qubic/software/qubic/job_rpc_server.py /home/qubic/jobserver_config.yaml
