# packages/flow-compute/flow_compute-0.0.20.tar.gz/flow_compute-0.0.20/src/flow/providers/mithril/runtime/startup/templates/sections/foundrypf.service.j2
[Unit]
Description=Foundry Port Forwarding
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/foundrypf
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foundrypf

[Install]
WantedBy=multi-user.target


